# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# PostgreSQL Database
[program:postgresql]
command=/usr/lib/postgresql/15/bin/postgres -D /data/postgresql
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/postgresql.log
stderr_logfile=/var/log/supervisor/postgresql_err.log

# Redis Server
[program:redis]
command=redis-server --dir /data/redis --appendonly yes --requirepass redis_secure_password
user=neuroinsight
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis_err.log

# MinIO Object Storage
[program:minio]
command=/usr/local/bin/minio server /data/minio --console-address ":9001"
user=neuroinsight
autostart=true
autorestart=true
priority=30
environment=MINIO_ROOT_USER="minioadmin",MINIO_ROOT_PASSWORD="minioadmin_secure"
stdout_logfile=/var/log/supervisor/minio.log
stderr_logfile=/var/log/supervisor/minio_err.log

# Backend API Server
[program:backend]
command=/usr/local/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000 --workers 2
directory=/app
user=neuroinsight
autostart=true
autorestart=true
priority=40
environment=PYTHONPATH="/app",POSTGRES_HOST="localhost",REDIS_HOST="localhost",MINIO_ENDPOINT="localhost:9000"
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend_err.log

# Celery Worker
[program:celery-worker]
command=/usr/local/bin/celery -A workers.tasks.processing_web worker --loglevel=info --concurrency=2
directory=/app
user=neuroinsight
autostart=true
autorestart=true
priority=50
environment=PYTHONPATH="/app",POSTGRES_HOST="localhost",REDIS_HOST="localhost",MINIO_ENDPOINT="localhost:9000"
stdout_logfile=/var/log/supervisor/celery_worker.log
stderr_logfile=/var/log/supervisor/celery_worker_err.log
stopwaitsecs=600
stopasgroup=true
killasgroup=true

# Celery Beat Scheduler
[program:celery-beat]
command=/usr/local/bin/celery -A workers.tasks.processing_web beat --loglevel=info
directory=/app
user=neuroinsight
autostart=true
autorestart=true
priority=60
environment=PYTHONPATH="/app",POSTGRES_HOST="localhost",REDIS_HOST="localhost"
stdout_logfile=/var/log/supervisor/celery_beat.log
stderr_logfile=/var/log/supervisor/celery_beat_err.log

# Job Monitor
[program:job-monitor]
command=/usr/local/bin/python -c "import sys; sys.path.insert(0, '/app'); from backend.services.job_monitor import JobMonitor; monitor = JobMonitor(); monitor.start_background_monitoring(); import time; [time.sleep(60) for _ in iter(int, 1)]"
directory=/app
user=neuroinsight
autostart=true
autorestart=true
priority=70
environment=PYTHONPATH="/app",POSTGRES_HOST="localhost",REDIS_HOST="localhost"
stdout_logfile=/var/log/supervisor/job_monitor.log
stderr_logfile=/var/log/supervisor/job_monitor_err.log
