# docker/allinone/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/data/logs/supervisord.log
pidfile=/tmp/supervisord.pid

[program:postgres]
command=/opt/allinone/start-pg.sh
user=neuroinsight
autorestart=true
priority=10
stdout_logfile=/data/logs/postgres.log
stderr_logfile=/data/logs/postgres.log

[program:redis]
command=redis-server --bind 127.0.0.1 --dir /data/redis --save "" --requirepass nirredis
user=neuroinsight
autorestart=true
priority=10
stdout_logfile=/data/logs/redis.log
stderr_logfile=/data/logs/redis.log

[program:minio]
command=minio server /data/minio --address 127.0.0.1:9000 --console-address 127.0.0.1:9001
user=neuroinsight
environment=MINIO_ROOT_USER="niradmin",MINIO_ROOT_PASSWORD="nirminio123"
autorestart=true
priority=10
stdout_logfile=/data/logs/minio.log
stderr_logfile=/data/logs/minio.log

[program:api]
command=/opt/allinone/start-api.sh
user=neuroinsight
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worker]
command=/opt/allinone/start-worker.sh
user=neuroinsight
autorestart=true
priority=30
stdout_logfile=/data/logs/worker.log
stderr_logfile=/data/logs/worker.log
