# systemd/neuroinsight-beat.service
[Unit]
Description=NeuroInsight Celery Beat Scheduler
After=network.target
Requires=neuroinsight-backend.service neuroinsight-worker.service

[Service]
Type=simple
WorkingDirectory=%h/src/desktop_alone_web_1
EnvironmentFile=%h/src/desktop_alone_web_1/.env
Environment="PYTHONPATH=%h/src/desktop_alone_web_1"
Environment="ENVIRONMENT=production"
ExecStart=%h/src/desktop_alone_web_1/venv/bin/python -m celery -A workers.tasks.processing_web beat --loglevel=info
StandardOutput=append:%h/src/desktop_alone_web_1/celery_beat.log
StandardError=append:%h/src/desktop_alone_web_1/celery_beat.log
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
