# systemd/neuroinsight-backend.service
[Unit]
Description=NeuroInsight Backend API Server
After=network.target
Wants=neuroinsight-worker.service neuroinsight-beat.service

[Service]
Type=simple
WorkingDirectory=%h/src/desktop_alone_web_1
EnvironmentFile=%h/src/desktop_alone_web_1/.env
Environment="PYTHONPATH=%h/src/desktop_alone_web_1"
Environment="ENVIRONMENT=production"
Environment="API_PORT=8000"
Environment="PORT=8000"
Environment="MAX_CONCURRENT_JOBS=1"
ExecStart=%h/src/desktop_alone_web_1/venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000
StandardOutput=append:%h/src/desktop_alone_web_1/neuroinsight.log
StandardError=append:%h/src/desktop_alone_web_1/neuroinsight.log
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
