# systemd/neuroinsight-monitor.service
[Unit]
Description=NeuroInsight Job Monitor
After=network.target
Requires=neuroinsight-backend.service

[Service]
Type=simple
WorkingDirectory=%h/src/desktop_alone_web_1
EnvironmentFile=%h/src/desktop_alone_web_1/.env
Environment="PYTHONPATH=%h/src/desktop_alone_web_1"
Environment="ENVIRONMENT=production"
Environment="FREESURFER_CONTAINER_PREFIX=freesurfer-job-"
ExecStart=%h/src/desktop_alone_web_1/venv/bin/python -c "import sys; sys.path.insert(0, '.'); from backend.services.job_monitor import JobMonitor; monitor = JobMonitor(); monitor.start_background_monitoring(); import time; [time.sleep(60) for _ in iter(int, 1)]"
StandardOutput=append:%h/src/desktop_alone_web_1/job_monitor.log
StandardError=append:%h/src/desktop_alone_web_1/job_monitor.log
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
